# btrfs-slots/files/etc/uci-defaults/09-fstab
[ -f /etc/config/fstab ] && exit 0

#write fstab before the default `block detect`, enable /efi
uci -q batch <<-EOF >/dev/null
    add fstab global
    set fstab.@global[-1].anon_swap=0
    set fstab.@global[-1].anon_mount=0
    set fstab.@global[-1].auto_swap=1
    set fstab.@global[-1].auto_mount=1
    set fstab.@global[-1].delay_root=5
    set fstab.@global[-1].check_fs=0

    add fstab mount
    set fstab.@mount[-1].target=/efi
    set fstab.@mount[-1].label=efi
    set fstab.@mount[-1].enabled=1
    set fstab.@mount[-1].fstype=vfat
    commit fstab
EOF
